# germinal.service
[Unit]
Description=Germinal, a gemini server

[Service]
ExecStart=/usr/local/share/germinal/germinal.ros
Type=simple
Restart=always
RestartSec=10
User=germinal
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes


[Install]
WantedBy=multi-user.target
